# scripts/99-gz302-rgb.rules
# GZ302 RGB Control - udev rules for unprivileged access
# This file configures permissions for both keyboard and rear window RGB

# ============================================================================
# Keyboard RGB Control - ASUS ROG Flow Z13 keyboard (USB 0b05:1a30)
# ============================================================================
# USB device access for gz302-rgb binary (libusb mode)
SUBSYSTEMS=="usb", ATTRS{idVendor}=="0b05", ATTRS{idProduct}=="1a30", MODE="0666"

# HIDRAW device access for keyboard RGB (hidraw mode)
KERNEL=="hidraw*", ATTRS{idVendor}=="0b05", ATTRS{idProduct}=="1a30", MODE="0666"

# ============================================================================
# Rear Window/Lightbar RGB Control - N-KEY Device (USB 0b05:18c6)
# ============================================================================
# USB device access
SUBSYSTEMS=="usb", ATTRS{idVendor}=="0b05", ATTRS{idProduct}=="18c6", MODE="0666"

# HIDRAW device access for lightbar control
KERNEL=="hidraw*", ATTRS{idVendor}=="0b05", ATTRS{idProduct}=="18c6", MODE="0666"

# ============================================================================
# Keyboard Backlight sysfs permissions (asus driver)
# ============================================================================
# Allows unprivileged write access to kbd_backlight brightness
ACTION=="add|change", SUBSYSTEM=="leds", KERNEL=="*::kbd_backlight", RUN+="/bin/chmod 0666 /sys/%p/brightness"
ACTION=="add|change", SUBSYSTEM=="leds", KERNEL=="asus::kbd_backlight", RUN+="/bin/chmod 0666 /sys/%p/brightness"

# ASUS aura_keyboard interface (kernel 6.14+)
ACTION=="add|change", SUBSYSTEM=="leds", KERNEL=="aura_keyboard", RUN+="/bin/chmod -R 0666 /sys/%p/*"
