# scripts/99-gz302-rgb.rules
# GZ302 RGB Control - udev rules for unprivileged access
# This file configures permissions for both keyboard and rear window RGB

# Keyboard RGB Control - ASUS ROG Flow Z13 keyboard (USB 0b05:1a30)
# Allows unprivileged access to the USB device for gz302-rgb binary
SUBSYSTEMS=="usb", ATTRS{idVendor}=="0b05", ATTRS{idProduct}=="1a30", MODE="0666"

# Rear Window/Lightbar RGB Control - sysfs LED brightness
# Allows unprivileged write access to kbd_backlight brightness files
# Multiple patterns to catch different device paths
ACTION=="add|change", SUBSYSTEM=="leds", KERNEL=="*::kbd_backlight", RUN+="/bin/chmod 0666 /sys/%p/brightness"
ACTION=="add|change", SUBSYSTEM=="leds", KERNEL=="asus::kbd_backlight", RUN+="/bin/chmod 0666 /sys/%p/brightness"

# Lightbar-specific LED device (if exposed as separate device)
SUBSYSTEMS=="usb", ATTRS{idVendor}=="0b05", ATTRS{idProduct}=="18c6", MODE="0666"
